# roles/api/templates/upstart/upstart.conf
description "{{ service }}"

start on started networking
stop on runlevel [016]

limit nofile 1000000 1000000

respawn
respawn limit 5 30

setuid {{ service }}
setgid {{ service }}

chdir /home/{{service}}/src
env NODE_ENV=production
env NODE_APP_INSTANCE={{node_app_instance}}
exec npm start
